# deploy/raspberry-pi/systemd/prevue-kiosk.service
[Unit]
Description=Prevue Kiosk Display (Epiphany)
Documentation=https://github.com/user/prevue
Requires=prevue-docker.service
After=prevue-docker.service graphical.target display-manager.service
PartOf=prevue.target

[Service]
Type=simple
User=prevue
Group=prevue
Environment="DISPLAY=:0"
Environment="HOME=/home/prevue"
ExecStart=/home/prevue/deploy/kiosk/start-kiosk.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=no
ProtectHome=no

[Install]
WantedBy=graphical.target
