# deploy/raspberry-pi/systemd/prevue-kiosk.service
[Unit]
Description=Prevue Kiosk Display (Epiphany on Pi OS Desktop)
Documentation=https://github.com/qrobinso/prevue
Requires=prevue-docker.service
After=prevue-docker.service graphical.target
PartOf=prevue.target

[Service]
Type=simple
User=prevue
Group=prevue
Environment="DISPLAY=:0"
Environment="HOME=/home/prevue"
Environment="XAUTHORITY=/home/prevue/.Xauthority"
ExecStart=/home/prevue/deploy/raspberry-pi/kiosk/start-kiosk.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=no
ProtectHome=no

[Install]
WantedBy=graphical.target
